# deploy/ficc-connect.service
[Unit]
Description=FICC Connect Backend API Server
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ficc-connect/backend
ExecStart=/usr/bin/node /opt/ficc-connect/backend/server.js
Restart=always
RestartSec=10
StandardOutput=append:/opt/ficc-connect/backend/logs/app.log
StandardError=append:/opt/ficc-connect/backend/logs/error.log

# Load environment variables from .env file
EnvironmentFile=/opt/ficc-connect/backend/.env

# Additional environment
Environment=NODE_ENV=production
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ficc-connect/backend/logs

[Install]
WantedBy=multi-user.target
